# debian/postinst
#!/bin/bash
set -e

# Source debconf library
. /usr/share/debconf/confmodule

# Source common installation functions
if [ -f /opt/joblet/scripts/common-install-functions.sh ]; then
    . /opt/joblet/scripts/common-install-functions.sh
else
    echo "ERROR: Common installation functions not found!"
    exit 1
fi

# Debian-specific configuration with debconf support
get_configuration_debian() {
    # Configuration precedence (highest to lowest):
    # 1. Environment variables (for automated deployments)
    # 2. Debconf (for interactive installations)
    # 3. Auto-detection (fallback)

    print_info "Loading configuration..."

    # === Priority 1: Environment Variables ===
    if [ -n "$JOBLET_SERVER_ADDRESS" ] || [ -n "$JOBLET_CERT_INTERNAL_IP" ] || [ -n "$JOBLET_SERVER_IP" ]; then
        print_info "Configuration source: Environment variables"

        # Legacy support: JOBLET_SERVER_IP sets the certificate primary IP
        if [ -n "$JOBLET_SERVER_IP" ]; then
            JOBLET_CERT_PRIMARY="${JOBLET_CERT_PRIMARY:-$JOBLET_SERVER_IP}"
            JOBLET_CERT_INTERNAL_IP="${JOBLET_CERT_INTERNAL_IP:-$JOBLET_SERVER_IP}"
            print_info "Using JOBLET_SERVER_IP (legacy): $JOBLET_SERVER_IP"
        fi

        JOBLET_SERVER_ADDRESS="${JOBLET_SERVER_ADDRESS:-0.0.0.0}"
        JOBLET_SERVER_PORT="${JOBLET_SERVER_PORT:-50051}"

    # === Priority 2: Debconf (Debian-specific) ===
    elif command -v db_get >/dev/null 2>&1; then
        print_info "Configuration source: Debconf (interactive)"

        db_get joblet/server_address || true
        JOBLET_SERVER_ADDRESS="${RET:-0.0.0.0}"

        db_get joblet/server_port || true
        JOBLET_SERVER_PORT="${RET:-50051}"

        db_get joblet/cert_internal_ip || true
        JOBLET_CERT_INTERNAL_IP="${RET}"

        db_get joblet/cert_public_ip || true
        JOBLET_CERT_PUBLIC_IP="${RET}"

        db_get joblet/cert_domain || true
        JOBLET_CERT_DOMAIN="${RET}"

    # === Priority 3: Use common function ===
    else
        # Call the common get_configuration function
        get_configuration
        return 0
    fi

    # === Auto-detect internal IP if not set (all paths) ===
    if [ -z "$JOBLET_CERT_INTERNAL_IP" ]; then
        JOBLET_CERT_INTERNAL_IP=$(detect_internal_ip)
        print_info "Auto-detected internal IP: $JOBLET_CERT_INTERNAL_IP"
    fi

    # === Set primary certificate address (used for CN) ===
    JOBLET_CERT_PRIMARY=${JOBLET_CERT_PRIMARY:-$JOBLET_CERT_INTERNAL_IP}

    # === Build SAN list for certificate ===
    if [ -z "$JOBLET_ADDITIONAL_NAMES" ]; then
        JOBLET_ADDITIONAL_NAMES="localhost"

        if [ -n "$JOBLET_CERT_INTERNAL_IP" ] && [ "$JOBLET_CERT_INTERNAL_IP" != "$JOBLET_CERT_PRIMARY" ]; then
            JOBLET_ADDITIONAL_NAMES="$JOBLET_ADDITIONAL_NAMES,$JOBLET_CERT_INTERNAL_IP"
        fi

        if [ -n "$JOBLET_CERT_PUBLIC_IP" ]; then
            JOBLET_ADDITIONAL_NAMES="$JOBLET_ADDITIONAL_NAMES,$JOBLET_CERT_PUBLIC_IP"
        fi

        if [ -n "$JOBLET_CERT_DOMAIN" ]; then
            JOBLET_ADDITIONAL_NAMES="$JOBLET_ADDITIONAL_NAMES,$JOBLET_CERT_DOMAIN"
        fi
    fi

    print_success "Configuration loaded successfully"
}

case "$1" in
    configure)
        # Display system changes warning
        display_system_changes_warning

        print_info "🔧 Configuring Joblet Service with Enhanced Certificate Support..."
        echo

        # Set basic permissions first
        chown -R root:root /opt/joblet
        chmod 755 /opt/joblet
        chmod 755 /opt/joblet/bin
        chmod 755 /opt/joblet/bin/joblet
        chmod 755 /opt/joblet/bin/rnx
        chmod 755 /opt/joblet/scripts
        chmod 644 /opt/joblet/scripts/joblet-config-template.yml
        chmod 644 /opt/joblet/scripts/rnx-config-template.yml
        chmod +x /usr/local/bin/certs_gen_embedded.sh

        mkdir -p /opt/joblet/config
        chmod 700 /opt/joblet/config  # Restricted since it will contain private keys

        # Create symlinks
        if [ ! -L /usr/bin/rnx ]; then
            ln -sf /opt/joblet/bin/rnx /usr/bin/rnx
        fi

        if [ ! -L /usr/local/bin/rnx ]; then
            ln -sf /opt/joblet/bin/rnx /usr/local/bin/rnx
        fi

        get_configuration_debian

        # Detect AWS environment and display information
        detect_aws_environment

        print_info "Configuration Summary:$EC2_INFO"
        echo "  gRPC Server Bind: $JOBLET_SERVER_ADDRESS:$JOBLET_SERVER_PORT"
        echo "  Certificate Primary IP: $JOBLET_CERT_PRIMARY"
        if [ -n "$JOBLET_CERT_PUBLIC_IP" ]; then
            echo "  Certificate Public IP: $JOBLET_CERT_PUBLIC_IP"
        fi
        if [ -n "$JOBLET_CERT_DOMAIN" ]; then
            echo "  Certificate Domain(s): $JOBLET_CERT_DOMAIN"
        fi
        echo

        # Generate certificates and embed them in config files
        if generate_and_embed_certificates; then
            # Set secure permissions on config files (they now contain private keys)
            chmod 600 /opt/joblet/config/joblet-config.yml
            chmod 600 /opt/joblet/config/rnx-config.yml

            # Create convenience copy for local CLI usage
            if [ -f /opt/joblet/config/rnx-config.yml ]; then
                mkdir -p /etc/joblet
                cp /opt/joblet/config/rnx-config.yml /etc/joblet/rnx-config.yml
                chmod 644 /etc/joblet/rnx-config.yml
                print_success "Client config available at /etc/joblet/rnx-config.yml"
            fi
        else
            print_error "Failed to generate certificates"
            exit 1
        fi

        # Network requirements configured
        setup_network_requirements

        # Create runtime directories
        mkdir -p /var/log/joblet
        mkdir -p /opt/joblet/logs
        mkdir -p /opt/joblet/network
        mkdir -p /opt/joblet/volumes
        mkdir -p /opt/joblet/jobs
        
        # Set ownership and permissions
        chown root:root /var/log/joblet /opt/joblet/logs /opt/joblet/network /opt/joblet/volumes /opt/joblet/jobs
        chmod 755 /var/log/joblet /opt/joblet/logs /opt/joblet/network /opt/joblet/volumes /opt/joblet/jobs

        # Setup cgroup delegation
        if [ -d /sys/fs/cgroup ]; then
            print_info "Setting up cgroup delegation..."
            mkdir -p /sys/fs/cgroup/joblet.slice
            echo "+cpu +memory +io +pids +cpuset" > /sys/fs/cgroup/joblet.slice/cgroup.subtree_control 2>/dev/null || true
        fi

        # Note: joblet-persist now runs as a subprocess of joblet (no separate service needed)
        print_info "joblet-persist will run as a subprocess of joblet-core"

        # Enable systemd services
        systemctl daemon-reload
        systemctl enable joblet.service

        # Clean up temporary files (EC2 info may be used for AWS deployments)
        rm -f /tmp/joblet-ec2-info

        # Display quickstart information
        display_quickstart_info "debian"

        # Setup user-local rnx installation
        if [ -n "$SUDO_USER" ] && [ "$SUDO_USER" != "root" ]; then
            USER_HOME=$(getent passwd "$SUDO_USER" | cut -d: -f6)
            if [ -n "$USER_HOME" ] && [ -d "$USER_HOME" ]; then
                print_info "Setting up rnx in user directory for $SUDO_USER..."

                # Create ~/.rnx directory
                mkdir -p "$USER_HOME/.rnx"

                # Copy rnx binary
                cp /opt/joblet/bin/rnx "$USER_HOME/.rnx/"
                chmod 755 "$USER_HOME/.rnx/rnx"

                # Copy config if it exists
                if [ -f /opt/joblet/config/rnx-config.yml ]; then
                    cp /opt/joblet/config/rnx-config.yml "$USER_HOME/.rnx/"
                    chmod 644 "$USER_HOME/.rnx/rnx-config.yml"
                fi

                # Set ownership
                chown -R "$SUDO_USER:$(id -gn "$SUDO_USER")" "$USER_HOME/.rnx"

                # Add to PATH in shell profiles if not already present
                for profile in "$USER_HOME/.bashrc" "$USER_HOME/.profile" "$USER_HOME/.zshrc"; do
                    if [ -f "$profile" ]; then
                        if ! grep -q 'export PATH="$HOME/.rnx:$PATH"' "$profile"; then
                            echo '' >> "$profile"
                            echo '# Add ~/.rnx to PATH for joblet rnx CLI' >> "$profile"
                            echo 'export PATH="$HOME/.rnx:$PATH"' >> "$profile"
                            chown "$SUDO_USER:$(id -gn "$SUDO_USER")" "$profile"
                            print_success "Added ~/.rnx to PATH in $profile"
                        fi
                    fi
                done

                print_success "rnx installed in $USER_HOME/.rnx"
                print_info "Reload your shell or run: source ~/.bashrc"
            fi
        fi

        ;;
esac

exit 0